# systemd/scribe.timer
[Unit]
Description=Weekly Timer for scribe

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
